# Ukenytt_Frida/rootfs/etc/s6-overlay/s6-rc.d/ukenytt/finish
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Ukenytt Add-on: Håndterer avslutning
# ==============================================================================

bashio::log.info "Ukenytt add-on avsluttet med kode: ${1}"

if [[ "${1}" -ne 0 ]] && [[ "${1}" -ne 256 ]]; then
    bashio::log.warning "Ukenytt krasjet, venter 5 sekunder før restart..."
    sleep 5
fi
