# Ukenytt_Frida/rootfs/etc/s6-overlay/s6-rc.d/ukenytt/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Ukenytt Add-on: Starter hovedapplikasjonen
# ==============================================================================

bashio::log.info "Starter Ukenytt add-on..."

# Opprett data-mappe hvis den ikke finnes
mkdir -p /data/ukenytt

# Hent konfigurasjon via bashio
export UKENYTT_API_KEY=$(bashio::config 'api_key')
# Les children som JSON-array direkte fra options.json
export UKENYTT_CHILDREN=$(jq -c '.children' /data/options.json)

# Bruk Home Assistant API proxy
export UKENYTT_HA_URL="http://supervisor/core"
export UKENYTT_HA_TOKEN="${SUPERVISOR_TOKEN}"

bashio::log.info "Konfigurerte barn: ${UKENYTT_CHILDREN}"
bashio::log.info "Starter webserver på port 8099..."

# Test Python-import først
bashio::log.info "Tester Python-avhengigheter..."
if ! python3 -c "import pandas; import tabula; import flask; import requests; print('OK')" 2>&1; then
    bashio::log.error "Python-import feilet!"
fi

# Start Python-applikasjonen med full feilutskrift
python3 -u /usr/bin/ukenytt.py 2>&1 || bashio::log.error "Python avsluttet med feil: $?"
exec sleep infinity
